# systemd/bungeecord.service
[Unit]
Description=Minecraft server - BungeeCord
After=local-fs.target network.target minecraft.service

[Service]
WorkingDirectory=/home/ubuntu/minecraft
User=ubuntu
Group=minecraft
Type=simple
Restart=on-failure
RestartSec=10
ExecStart=/home/ubuntu/minecraft/bin/bungee-start.sh

[Install]
WantedBy=multi-user.target